# launchd/jobs/com.onyx.claude.job.evolve-cycle.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.onyx.claude.job.evolve-cycle</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__HOME__/__INSTALL_DIR__/scripts/run-job.sh</string>
        <string>evolve-cycle</string>
        <string>--notify</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
        
    </dict>
    <key>StandardOutPath</key>
    <string>__HOME__/onyx-claude-logs/evolve-cycle-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/onyx-claude-logs/evolve-cycle-launchd-err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
</dict>
</plist>
